# voice_bot_vk.service
[Service]
ExecStart=/opt/voice_bot/.venv/bin/python3 /opt/voice_bot/run_vk_bot.py
Restart=always

[Install]
WantedBy=multi-user.target
